# templates/docker-cleanup.timer.j2
[Unit]
Description=Docker cleanup timer

[Timer]
OnUnitInactiveSec={{ docker_extra_cleanup_interval }}

[Install]
WantedBy=timers.target